# packages/melopero-autostart/melopero_autostart-0.1.0.tar.gz/melopero_autostart-0.1.0/melopero-autostart/melopero-autostart.service
[Unit]
Description=melopero-autostart
After=multi-user.target

[Service]
Type=forking
ExecStart=/home/melopero-autostart/StartScripts.sh

[Install]
WantedBy=multi-user.target
