# packages/libreeye/libreeye-0.4.1-py3-none-any.whl/libreeye/package_data/etc/systemd/libreeye.service
[Unit]
Description=Libreeye daemon service
After=docker.service

[Service]
Type=simple
PIDFile=/run/libreeye.pid
ExecStart=/usr/bin/env python3 -m libreeye.daemon.daemon
ExecStop=/bin/kill -s SIGTERM $MAINPID
TimeoutStopSec=infinity
User=root

[Install]
WantedBy=multi-user.target