# packages/libreeye/libreeye-0.3.0.tar.gz/libreeye-0.3.0/src/libreeye/package_data/etc/systemd/libreeye.service
[Unit]
Description=Libreeye daemon service

[Service]
Type=forking
PIDFile=/run/libreeye.pid
ExecStart=/usr/bin/env python3 -m libreeye.daemon.daemon
ExecStop=/bin/kill -s SIGTERM $MAINPID
User=root

[Install]
WantedBy=multi-user.target